# etc/inittab
#
tts/0::sysinit:/etc/init.d/rcS > /dev/console 2> /dev/console

# ARRIS CHANGE - Spawn the ARRIS login shell instead of ash by default
tts/0::respawn:/usr/sbin/mini_cli

# Stuff to do before rebooting
#::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -l /nvram /var /sys
null::shutdown:ifconfig usb0 down 2> /dev/null
#::shutdown:/sbin/swapoff -a
